# system-maintenance-suite/etc/systemd/maintenance-backup.timer
[Unit]
Description=Run maintenance backup hourly

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
