# system-maintenance-suite/etc/systemd/maintenance-backup.service
[Unit]
Description=Maintenance suite backup

[Service]
Type=oneshot
EnvironmentFile=/etc/maintenance-suite/env
ExecStart=/usr/local/bin/backup.sh
